# provisioning/roles/ipfs-tika/templates/ipfs-tika.service
[Unit]
Description=IPFS Tika daemon
Requires=ipfs.service
After=ipfs.service

[Service]
ExecStart=/usr/bin/java -Xmx256m -jar /usr/local/libexec/{{ ipfs_tika_jar }}
User=ipfs-tika
Group=ipfs-tika
MemoryAccounting=yes
MemoryLimit=512M

[Install]
WantedBy=default.target
